# scripts/crontab
15   0   *   *   0   %%USER_NAME%%    cd %%INSTALL_DIR%%; ./backup.sh
