# scripts/crontab
15   0   *   *   0   %%USER_NAME%%    cd %%BIN_DIR%%; ./backup.sh
